# ui/ai.openclaw.context-ui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.context-ui</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <!-- UPDATE: path to context-ui.py -->
        <string>/path/to/openclaw-dynamic-context/ui/context-ui.py</string>
    </array>
    <!-- UPDATE: your workspace root (where context.yaml lives) -->
    <key>WorkingDirectory</key>
    <string>/path/to/workspace</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!-- UPDATE: log paths -->
    <key>StandardOutPath</key>
    <string>/path/to/logs/context-ui.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/logs/context-ui.err</string>
</dict>
</plist>
